# roles/common_linux/files/scripts/rc.local
#!/bin/bash

touch /var/lock/subsys/local

/etc/se_boot_all.sh

